# etc/udev/rules.d/91-backlight.rules
# Storing brightness
ACTION=="change", SUBSYSTEM=="backlight", RUN+="/bin/sh /etc/backlight.sh"
